# deploy/dialog-bot.service
# Copyright (C) 2026 Placeholder Name
# SPDX-License-Identifier: GPL-3.0-only

[Unit]
Description=Dialog Discord Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/dialog-bot
EnvironmentFile=/opt/dialog-bot/.env
ExecStart=/opt/dialog-bot/.venv/bin/python /opt/dialog-bot/bot.py
Restart=always
RestartSec=10
User=bitnami

[Install]
WantedBy=multi-user.target
